# ops/local/launchd/org.specdojo.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>org.specdojo.ollama</string>

        <key>ProgramArguments</key>
        <array>
            <string>/opt/homebrew/bin/ollama</string>
            <string>serve</string>
        </array>

        <key>EnvironmentVariables</key>
        <dict>
            <key>OLLAMA_KEEP_ALIVE</key>
            <string>10m</string>
            <key>OLLAMA_MAX_LOADED_MODELS</key>
            <string>2</string>
        </dict>

        <key>RunAtLoad</key>
        <true />

        <key>KeepAlive</key>
        <true />

        <key>StandardOutPath</key>
        <string>/tmp/ollama.out.log</string>

        <key>StandardErrorPath</key>
        <string>/tmp/ollama.err.log</string>
    </dict>
</plist>
